# client/updater.cron.sample
#
# Cron file for DynamicDNS client script
#
0 * * * * /usr/bin/php -q /path/to/your/updater.php
